# edge/scripts/aura-watch-edge.service.template
[Unit]
Description=Aura Watch AI Edge Surveillance Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
Environment=PYTHONUNBUFFERED=1
ExecStart=__PYTHON__ -u __DIR__/main.py
# Restart unless explicitly stopped via `systemctl stop`
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
