# edge/scripts/aura-watch-edge.service.template
[Unit]
Description=Aura Watch AI Edge Surveillance Agent
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__DIR__
ExecStart=__PYTHON__ __DIR__/main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
